# scripts/launchd/com.tom.me.review-decisions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tom.me.review-decisions</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>source ~/.zshrc &gt;/dev/null 2&gt;&amp;1 || true; exec pnpm review-decisions</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <!-- review-decisions is a pure deterministic scanner; no ANTHROPIC_API_KEY needed. -->
    <!-- If you want to set environment variables, add an EnvironmentVariables dict here. -->

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tom/Library/Logs/me/review-decisions.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tom/Library/Logs/me/review-decisions.err.log</string>
</dict>
</plist>
